# promtail/rootfs/etc/services.d/promtail/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Promtail
# Runs Promtail
# ==============================================================================

bashio::log.info 'Starting Promtail...'
promtail_config='/etc/promtail/config.yaml'

journal_path='/var/log/journal'
if ! bashio::fs.directory_exists "${journal_path}" || [ -z "$(ls -A ${journal_path})" ]; then
    bashio::log.info "No journal at ${journal_path}, looking for journal in /run/log/journal instead"
    journal_path='/run/log/journal'
fi

case "$(bashio::config 'log_level')" in \
    trace)      ;& \
    debug)      log_level='debug' ;; \
    notice)     ;& \
    warning)    log_level='warn' ;; \
    error)      ;& \
    fatal)      log_level='error' ;; \
    *)          log_level='info' ;; \
esac;
bashio::log.info "Promtail log level set to ${log_level}"

export "URL=$(bashio::config 'client.url')"
export "JOURNAL_PATH=${journal_path}"
export "LOG_LEVEL=$(bashio::config 'log_level')"

promtail_args=("-config.expand-env=true" "-config.file=${promtail_config}")
if [ "${log_level}" == "debug" ]; then
    bashio::log.debug "Logging full config on startup for debugging..."
    promtail_args+=("-print-config-stderr=true")
fi

bashio::log.info "Handing over control to Promtail..."
promtail "${promtail_args[@]}"
